# root/etc/crontabs/abc
# min   hour    day     month   weekday command
*       *       *       *       *       /usr/bin/php /app/www/public/crons/stats.php
*/5     *       *       *       *       /usr/bin/php /app/www/public/crons/state.php
*/5     *       *       *       *       /usr/bin/php /app/www/public/crons/pulls.php
*/10    *       *       *       *       /usr/bin/php /app/www/public/crons/housekeeper.php
*/15    *       *       *       *       /usr/bin/php /app/www/public/crons/health.php
0       12      *       *       *       /usr/bin/php /app/www/public/crons/prune.php