# root/etc/crontabs/abc
*/5 * * * * /usr/bin/php /config/www/crons/state.php
*/60 * * * * /usr/bin/php /config/www/crons/pulls.php