# root/etc/s6-overlay/s6-rc.d/init-dockwatch-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo 's6-overlay init require(/app/www/public/startup.php) ->';
s6-setuidgid abc php -r "require '/app/www/public/startup.php';"
echo 's6-overlay init require(/app/www/public/startup.php) <-';

# permissions
find /config ! -user abc -exec chown abc {} +
