# root/etc/s6-overlay/s6-rc.d/init-dockwatch-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

php -r "require '/app/www/public/startup.php';" 2>/dev/null
