# root/etc/s6-overlay/s6-rc.d/init-dockwatch-config/up
/etc/s6-overlay/s6-rc.d/init-dockwatch-config/run
