# root/etc/s6-overlay/s6-rc.d/init-memcache-server/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

s6-setuidgid abc memcached -u abc -m 64 -l 127.0.0.1
